# systemd/claude-telegram-bot.service
[Unit]
Description=Claude Telegram Bot
After=network.target

[Service]
Type=simple
User=claudebot
WorkingDirectory=/home/claudebot/claude-telegram-vps-bot
Environment=HOME=/home/claudebot
Environment=PATH=/home/claudebot/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/claudebot/claude-telegram-vps-bot/.env
ExecStart=/home/claudebot/claude-telegram-vps-bot/venv/bin/python /home/claudebot/claude-telegram-vps-bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
