# systemd/claude-telegram-bot.service
[Unit]
Description=Claude Telegram Bot
After=network.target

[Service]
Type=simple
User=tsoillc
WorkingDirectory=/home/tsoillc/claude-telegram-vps-bot
Environment=HOME=/home/tsoillc
Environment=PATH=/home/tsoillc/.npm-global/bin:/home/tsoillc/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/tsoillc/claude-telegram-vps-bot/.env
ExecStart=/home/tsoillc/claude-telegram-vps-bot/venv/bin/python /home/tsoillc/claude-telegram-vps-bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
